# deploy/aksystems-web.service
# Альтернатива pm2 — systemd-юнит.
# /etc/systemd/system/aksystems-web.service
# sudo systemctl enable --now aksystems-web

[Unit]
Description=AK Systems — Next.js production server
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/aksystems/web
EnvironmentFile=/var/www/aksystems/web/.env.local
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /var/www/aksystems/web/node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
